# devops/cloud-server/etc/cron.d/incanto-backup
01 01 * * * root /bin/bash /opt/backup.sh >> /var/log/cronjob.log
01 04 * * * root /bin/bash /opt/bkp-rotate.sh >> /var/log/bkp-rotate.log
