# devops/cloud-server/home/root/crontab
05,15,25,35,45,55 * * * * /opt/inbox-files.sh 2>&1 /var/log/inbox-files.log